# watch-config.path
[Unit]
Description="Monitor config.txt for configuration changes"

[Path]
PathChanged=/boot/config.txt
Unit=watch-config.service

[Install]
WantedBy=multi-user.target
